# com.ailearning.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ailearning.automation</string>
    
    <key>Program</key>
    <string>/usr/local/bin/node</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/arog/ai-lessons/daemon.js</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/arog/ai-lessons</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/arog/ai-lessons/logs/daemon.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/arog/ai-lessons/logs/daemon.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>